# deploy/obey-bot.service
[Unit]
Description=Obey RP Discord bot (announcements <-> forum bridge)
After=network.target obey-bridge.service
Wants=network-online.target obey-bridge.service

[Service]
Type=simple
User=obey
Group=obey
WorkingDirectory=/home/obey/obey-bridge
EnvironmentFile=/home/obey/obey-bridge/.env
ExecStart=/usr/bin/node /home/obey/obey-bridge/dist/bot/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening — same posture as obey-bridge.service. Bot only opens an
# outbound Discord WS + loopback HTTP to the bridge; no inbound ports.
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome=true masks /home entirely, which makes the WorkingDirectory
# below unreachable (systemd reports 200/CHDIR on start). Match the
# runbook's bridge unit: read-only lets the service load its own files.
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/obey/obey-bridge
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
